# openwrt/files/etc/uci-defaults/99-cerberus-packages
#!/bin/sh
# Cerberus first-boot setup

LOG="/tmp/cerberus-firstboot.log"
echo "$(date) - Cerberus first-boot starting" > "$LOG"

# ── Set root password to toor (direct shadow write) ──────────
sed -i 's|^root:[^:]*:|root:$1$cerberus$WRKWzVSUTSNoD..XxJosQ1:|' /etc/shadow
echo "$(date) - Root password set to toor" >> "$LOG"

# ── Create working directories ───────────────────────────────
mkdir -p /tmp/cerberus
mkdir -p /tmp/cerberus-captures

# ── Install extra packages if network available ──────────────
TRIES=0
while [ $TRIES -lt 12 ]; do
    if ping -c 1 -W 2 8.8.8.8 > /dev/null 2>&1; then
        echo "$(date) - Network available" >> "$LOG"
        opkg update >> "$LOG" 2>&1
        opkg install dsniff >> "$LOG" 2>&1 || true
        opkg install python3 python3-pip >> "$LOG" 2>&1 || true
        pip3 install sslstrip >> "$LOG" 2>&1 || true
        echo "$(date) - Packages done" >> "$LOG"
        break
    fi
    TRIES=$((TRIES + 1))
    sleep 5
done

echo "$(date) - First-boot done" >> "$LOG"
exit 0
