# openwrt/files/etc/uci-defaults/99-cerberus-packages
#!/bin/sh
# Cerberus first-boot package installer
# Runs once on first boot, then self-deletes (uci-defaults behavior)

LOG="/tmp/cerberus-firstboot.log"
echo "$(date) - Cerberus first-boot installer starting" > "$LOG"

# Create working directories
mkdir -p /tmp/cerberus
mkdir -p /tmp/cerberus-captures

# Try to install extra packages if we have internet
# This is best-effort — if no network, user can run manually later:
#   opkg update && opkg install dsniff python3 python3-pip && pip3 install sslstrip

# Wait for network (up to 60 seconds)
TRIES=0
while [ $TRIES -lt 12 ]; do
    if ping -c 1 -W 2 8.8.8.8 > /dev/null 2>&1; then
        echo "$(date) - Network available" >> "$LOG"
        opkg update >> "$LOG" 2>&1
        opkg install dsniff >> "$LOG" 2>&1 || true
        opkg install python3 python3-pip >> "$LOG" 2>&1 || true
        pip3 install sslstrip >> "$LOG" 2>&1 || true
        echo "$(date) - Package install complete" >> "$LOG"
        break
    fi
    TRIES=$((TRIES + 1))
    sleep 5
done

if [ $TRIES -ge 12 ]; then
    echo "$(date) - No network, skipping package install" >> "$LOG"
fi

echo "$(date) - First-boot script done" >> "$LOG"
exit 0
