# firmware/files/etc/uci-defaults/99-cerberus
#!/bin/sh
echo -e "toor\ntoor" | passwd root
uci set system.@system[0].hostname='Cerberus'
uci commit system
echo 1 > /proc/sys/net/ipv4/ip_forward
/etc/init.d/cerberus enable 2>/dev/null
exit 0
