# firmware/package/cerberus/files/cerberus.init
#!/bin/sh /etc/rc.common
START=99
STOP=10

start() {
    echo "Starting Cerberus..."
    mkdir -p /var/cerberus/captures /var/log
    echo 1 > /proc/sys/net/ipv4/ip_forward
    /usr/bin/cerberus > /var/log/cerberus.log 2>&1 &
}

stop() {
    echo "Stopping Cerberus..."
    killall cerberus 2>/dev/null
}

restart() {
    stop
    sleep 1
    start
}
